# provisioning/systemd/dobby-dbdump.service
[Unit]
Description=dobby nightly encrypted audit-DB dump
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/dobby
ExecStart=/opt/dobby/provisioning/dbdump.sh
