# provisioning/systemd/dobby-archive.service
[Unit]
Description=dobby monthly archive age-out
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/dobby
ExecStart=/usr/bin/docker compose run --rm archive-job run
# Monthly audit retention GC — prunes only sealed+expired rows and refuses to
# run if chain-verify finds a mismatch (never launders tampering).
ExecStart=/usr/bin/docker compose run --rm archive-job chain-gc
